# borg/.config/systemd/user/borg-backup.service
[Unit]
Description=Borg Backup Service

[Service]
Type=simple
ExecStart=/home/%u/bin/borg-backup.sh