# systemd/.config/systemd/user/borg-backup.service
[Unit]
Description=Borg Backup Service

[Service]
Type=simple
ExecStart=/home/andreas/bin/borg-backup.sh