# systemd/.config/systemd/user/sendip.timer
[Unit]
Description=Send ip to server every day

[Timer]
OnCalendar=daily
Persistent=true
Unit=sendip.service

[Install]
WantedBy=timers.target