# systemd/.config/systemd/user/sendip.service
[Unit]
Description=Send ip to server

[Service]
Type=simple
ExecStart=/home/andreas/bin/sendip.sh